# deploy/supervisor/thales-lb.conf
; Thales Transparent Load Balancer — supervisor config
; Deploy to: /etc/supervisor/conf.d/thales-lb.conf
; Reload: supervisorctl reread && supervisorctl update

; ─────────────────────────────────────────────────────────────────────────────
; Component 1: EzNet Inbound TCP→JMS adapter
; Listens on port 9100 for raw Thales commands from client applications.
; ─────────────────────────────────────────────────────────────────────────────
[program:thales-lb-inbound]
command=/usr/bin/java
    -Xms128m -Xmx384m
    -Djava.awt.headless=true
    -Deznet.root=/data1/xenticate/hsm-lb/eznet-inbound
    -Deznet.config.path=/data1/xenticate/hsm-lb/eznet-inbound/config
    -Deznet.log.path=/var/log/xenticate
    -jar /data1/xenticate/hsm-lb/bin/eznet-tcp2jms.war
directory=/data1/xenticate/hsm-lb/eznet-inbound
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/xenticate/thales-lb-inbound.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/xenticate/thales-lb-inbound-err.log
stderr_logfile_maxbytes=10MB
user=xenticate

; ─────────────────────────────────────────────────────────────────────────────
; Component 2: Thales Transparent Load Balancer (this project)
; Consumes JMS requests, routes raw Thales bytes to N HSM nodes.
; ─────────────────────────────────────────────────────────────────────────────
[program:thales-lb]
command=/usr/lib/jvm/temurin-25-jdk-arm64/bin/java
    -Xms256m -Xmx512m
    -Djava.awt.headless=true
    -jar /data1/xenticate/hsm-lb/bin/thales-transparent-lb.jar
    --spring.config.location=/data1/xenticate/hsm-lb/config/application.properties
directory=/data1/xenticate/hsm-lb
autostart=true
autorestart=true
startretries=5
startsecs=15
stopwaitsecs=30
stdout_logfile=/var/log/xenticate/thales-lb.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/xenticate/thales-lb-err.log
stderr_logfile_maxbytes=10MB
user=xenticate
