# deploy/supervisor/thales-lb.conf
; =============================================================================
; HSM Transparent Load Balancer — Supervisor Configuration
; =============================================================================
; Deploy to : /etc/supervisor/conf.d/thales-lb.conf
; Reload    : supervisorctl reread && supervisorctl update
; Start     : supervisorctl start thales-lb-inbound thales-lb
; Stop      : supervisorctl stop thales-lb thales-lb-inbound
; Status    : supervisorctl status thales-lb-inbound thales-lb
; =============================================================================

; -----------------------------------------------------------------------------
; Component 1: EzNet Inbound — TCP → JMS Adapter
;
; Listens on TCP port 9100 for raw HSM commands from client applications.
; Reads 2-byte length-prefixed frames and publishes to ActiveMQ queue
; hsm.transparent.lb.in. Sends replies back to client via JMS correlation.
; -----------------------------------------------------------------------------
[program:thales-lb-inbound]
command=/usr/lib/jvm/temurin-25-jdk-arm64/bin/java
    -Xms128m -Xmx384m
    -Djava.awt.headless=true
    -Dsun.net.http.allowRestrictedHeaders=true
    -Djava.io.tmpdir=/data1/xenticate/hsm-lb/eznet-inbound/tmp
    -Deznet.root=/data1/xenticate/hsm-lb/eznet-inbound
    -Deznet.config.path=/data1/xenticate/hsm-lb/eznet-inbound/config
    -Deznet.log.path=/var/log/xenticate
    -jar /data1/xenticate/hsm-lb/bin/eznet-tcp2jms.war

environment=JAVA_HOME="/usr/lib/jvm/temurin-25-jdk-arm64"
directory=/data1/xenticate/hsm-lb/eznet-inbound
user=xenticate
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
redirect_stderr=false
stdout_logfile=/var/log/xenticate/thales-lb-inbound.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/xenticate/thales-lb-inbound-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5

; -----------------------------------------------------------------------------
; Component 2: HSM Transparent Load Balancer
;
; Consumes JMS requests from hsm.transparent.lb.in, selects a healthy Thales
; node via configured algorithm, forwards raw bytes, returns raw response.
; Manages Commons Pool2 socket pools per node. Health checks every 5s.
; REST status: http://localhost:8110/api/v1/hsm-lb/status
; -----------------------------------------------------------------------------
[program:thales-lb]
command=/usr/lib/jvm/temurin-25-jdk-arm64/bin/java
    -Xms256m -Xmx512m
    -Djava.awt.headless=true
    -Dlog4j2.configurationFile=/data1/xenticate/hsm-lb/config/log4j2.xml
    -jar /data1/xenticate/hsm-lb/bin/thales-transparent-lb.jar
    --spring.config.location=/data1/xenticate/hsm-lb/config/application.properties

environment=JAVA_HOME="/usr/lib/jvm/temurin-25-jdk-arm64"
directory=/data1/xenticate/hsm-lb
user=xenticate
autostart=true
autorestart=true
startretries=5
startsecs=15
stopwaitsecs=30
redirect_stderr=false
stdout_logfile=/var/log/xenticate/thales-lb.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/xenticate/thales-lb-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
